# systemd/user/ai-agent.service
[Unit]
Description=AI Queue Agent (meta + claims + merge)
After=default.target

[Service]
Type=oneshot
Environment=XAIO_PROJECT_DIR=%h/xAIO-URL-Agent
EnvironmentFile=-%h/.config/xaio-url-agent.env
EnvironmentFile=-${XAIO_PROJECT_DIR}/.runtime/xaio-url-agent.env
WorkingDirectory=${XAIO_PROJECT_DIR}
ExecStart=/usr/bin/flock -n ./locks/ai.lock ${XAIO_PROJECT_DIR}/venv/bin/python ${XAIO_PROJECT_DIR}/src/ai_queue.py
