# systemd/user/pipeline.service
[Unit]
Description=xAIO pipeline

[Service]
Type=oneshot
EnvironmentFile=-%h/xAIO-URL-Agent/.runtime/.env
WorkingDirectory=%h/xAIO-URL-Agent
ExecStart=%h/xAIO-URL-Agent/venv/bin/python %h/xAIO-URL-Agent/src/pipeline_run.py --config %h/xAIO-URL-Agent/.runtime/config.yaml
