# systemd/user/pipeline.service
[Unit]
Description=xAIO pipeline

[Service]
Type=oneshot
Environment=XAIO_PROJECT_DIR=%h/xAIO-URL-Agent
EnvironmentFile=-%h/.config/xaio-url-agent.env
EnvironmentFile=-${XAIO_PROJECT_DIR}/.runtime/xaio-url-agent.env
WorkingDirectory=${XAIO_PROJECT_DIR}
ExecStart=${XAIO_PROJECT_DIR}/venv/bin/python ${XAIO_PROJECT_DIR}/src/pipeline_run.py
