# systemd/user/url-agent.service
[Unit]
Description=URL Agent (poll sheet, fetch pages, save JSON)
Wants=brave-agent.service
After=brave-agent.service

[Service]
Type=oneshot
Environment=XAIO_PROJECT_DIR=%h/xAIO-URL-Agent
EnvironmentFile=-%h/.config/xaio-url-agent.env
EnvironmentFile=-${XAIO_PROJECT_DIR}/.runtime/xaio-url-agent.env
WorkingDirectory=${XAIO_PROJECT_DIR}
ExecStart=/usr/bin/flock -n ./locks/capture.lock ${XAIO_PROJECT_DIR}/venv/bin/python ${XAIO_PROJECT_DIR}/src/agent.py
